# etc/local.d/20-setup-gpio.start
#!/bin/sh
chown root:gpio /dev/gpiomem
chmod g+rw /dev/gpiomem
insmod /opt/kmodules/devmem-full-access.ko || true
